# deploy/imtera-parser.service
[Unit]
Description=Imtera Yandex Maps Parser
After=network.target

[Service]
Type=simple
User=kp
WorkingDirectory=/var/www/imtera/services/parser
Environment=PORT=3001
ExecStart=/home/kp/.bun/bin/bun run start
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
